# vps/yscache_update.timer
[Unit]
Description=Daily update of UgraSage cache

[Timer]
OnCalendar=daily
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
